# deploy/systemd/quantgrid-backend.service
[Unit]
Description=QuantGrid backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/QuantGrid/services/trading-service
EnvironmentFile=/root/QuantGrid/services/trading-service/.env
ExecStart=/root/venv/bin/uvicorn Backend.presentation.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
